# deploy/filebrowser.service
[Unit]
Description=File Browser
After=network.target

[Service]
Type=simple
User=FILEBROWSER_USER
WorkingDirectory=FILEBROWSER_DIR
ExecStart=FILEBROWSER_DIR/.venv/bin/uvicorn filebrowser.main:app --host 127.0.0.1 --port FILEBROWSER_PORT
Restart=always
RestartSec=5
EnvironmentFile=FILEBROWSER_DIR/filebrowser.env

[Install]
WantedBy=multi-user.target
