# deploy/tailscale-cert-renew.timer
[Unit]
Description=Weekly Tailscale cert renewal

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
